# g2g/base/netopeer-cli/misc/supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface


[program:routes]
command=/bin/bash -c "chmod 755 routes.sh; ./routes.sh"
redirect_stderr=true
stdout_logfile=/home/netconf/cfgipsec2/netopeer.log
priority=2
autostart= True
autorestart = false
startsecs = 1
startretries = 0


